# patches/etc/init/init-virtualmin.conf
# init virtualmin configuration

description "initializes virtualmin on first start"
author "Scaleway <opensource@scaleway.com>"

start on (startup)

task

script
	/usr/local/sbin/oc-init-virtualmin
end script
