# patches/etc/init/init-webmin.conf
# init webmin configuration

description "initializes webmin on first start"
author "Scaleway <opensource@scaleway.com>"

start on (startup)

task

script
	/usr/local/sbin/oc-init-webmin
end script
